# docker/images/docker-spark/master/crontab
* 1 * * * /clueso/bin/table-compactor.sh /clueso/conf/application.conf 4 > /proc/1/fd/1 2>/proc/1/fd/2
