# docker/images/docker-spark/master/crontab
0 20 * * * root /clueso/bin/table-compactor.sh /clueso/conf/application.conf local 4
